# .config/systemd/user/udiskie.service
[Unit]
Description=Removable disk automounter using udisks

[Service]
ExecStart=/usr/bin/udiskie -ant2

[Install]
WantedBy=default.target
